# systemd/system/snap-gnome\x2d46\x2d2404-145.mount
[Unit]
Description=Mount unit for gnome-46-2404, revision 145
After=snapd.mounts-pre.target
Before=snapd.mounts.target

[Mount]
What=/var/lib/snapd/snaps/gnome-46-2404_145.snap
Where=/snap/gnome-46-2404/145
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=snapd.mounts.target
WantedBy=multi-user.target
